# arch/systemctl/rclone/rclone-backup.service
[Unit]
Description=Backup para OneDrive com rclone

[Service]
Type=oneshot
ExecStart=%h/scripts/rclone-backup.sh